# app/etc/cron/crontab
* * * * * root . $HOME/.profile; /usr/local/bin/php /app/bin/magento cron:run >> /app/var/log/magento.cron.log 2>&1